# com.gristmillcapital.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gristmillcapital.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kevinheaney/Desktop/Claude_Programs/Trading_Programs/gmc_dashboard/start_dashboard.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kevinheaney/Desktop/Claude_Programs/Trading_Programs/gmc_dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kevinheaney/Desktop/Claude_Programs/Trading_Programs/gmc_dashboard/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kevinheaney/Desktop/Claude_Programs/Trading_Programs/gmc_dashboard/dashboard_error.log</string>
</dict>
</plist>
